# tempout/tempout.service
[Unit]
Description=Outdoor temperature sensor

[Service]
ExecStart=/home/chip/git/k34/tempout/tempout.sh

[Install]
WantedBy=multi-user.target
