# tempout/tempgarage.service
[Unit]
Description=Garage temperature sensor

[Service]
ExecStart=/home/chip/git/k34/tempout/tempgarage.sh

[Install]
WantedBy=multi-user.target
